# deploy/discord_bot.service
[Unit]
Description=Discord Bot with automated pre-flight checks
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/discord_bot
# Optionally load environment variables (DISCORD_TOKEN, BOT_VENV_PATH, etc.)
# EnvironmentFile=/etc/discord_bot.env
ExecStart=/opt/discord_bot/scripts/deploy_update.sh
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
