# deploy/discord_bot.service
[Unit]
Description=HippoBot Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mars
WorkingDirectory=/home/mars/projects/discord_bot
Environment="PROJECT_DIR=/home/mars/projects/discord_bot"
ExecStart=/home/mars/projects/discord_bot/scripts/run_bot.sh
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hippobot

[Install]
WantedBy=multi-user.target
